# containers/cron-generate-daily-rss-dumps/crontab
# m h dom mon dow user    command
6   1 *   *   *   root    /opt/mediacloud/bin/cron_generate_daily_rss_dumps.pl /var/lib/daily_rss_dumps/ 1> /dev/stdout 2> /dev/stderr
