# apps/cron-set-media-primary-language/crontab
# Don't send no email
MAILTO=""

# m h dom mon dow user          command
16  2 *   *   *   mediacloud    /opt/mediacloud/bin/cron_set_media_primary_language.pl 1> /dev/stdout 2> /dev/stderr
