# apps/munin-cron/crontab
# Don't send no email
MAILTO=""

# m h dom mon dow user          command
*/5 * *   *   *   root          sudo -u munin /usr/bin/munin-cron --debug 1> /proc/1/fd/1 2> /proc/1/fd/2
