# apps/munin-cron/crontab
# m h dom mon dow user          command
*/5 * *   *   *   root          /munin-cron.sh 1> /dev/stdout 2> /dev/stderr
