# apps/munin-cron/crontab
# m h dom mon dow user          command
*/5 * *   *   *   munin         /usr/bin/munin-cron --debug 1> /dev/stdout 2> /dev/stderr
