# containers/cron-set-media-primary-language/crontab
# m h dom mon dow user    command
16  2 *   *   *   root    /cron_set_media_primary_language.pl 1> /proc/1/fd/1 2> /proc/1/fd/2
