# script/upstart_scripts/mcserver
#!/bin/sh
### BEGIN INIT INFO
# Provides:          mcserver
# Required-Start:    $syslog $remote_fs $network
# Required-Stop:     $syslog $remote_fs $network
# Should-Start:      mcserver
# Should-Stop:       mcserver
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start the Media Cloud web server.
### END INIT INFO


# Configuration
PATHTOMC=/home/mediacloud/trunk
DAEMONUSER=mediacloud

PATH=/sbin:/bin:/usr/sbin:/usr/bin
NAME=mcserver
DAEMON="$PATHTOMC/script/start_mediacloud_server.sh"
DESC="Media Cloud web server"
PIDFILE=/var/run/$NAME.pid
SCRIPTNAME=/etc/init.d/$NAME
DAEMON_ARGS=
LOGFILE=/var/log/mediacloud/$NAME.log

test -x $DAEMON || exit 0

set -e

if [ "$1" != status ]; then
    # be sure there is a /var/run/mediacloud, even with tmpfs
    mkdir --mode 750 --parents /var/run/mediacloud
    chown $DAEMONUSER /var/run/mediacloud

    # be sure there is a /var/log/mediacloud, even with tmpfs
    mkdir --mode 750 --parents /var/log/mediacloud
    chown $DAEMONUSER /var/log/mediacloud
fi

. /lib/lsb/init-functions

case "$1" in
    start)
        log_daemon_msg "Starting $DESC" "$NAME"
        cd "$PATHTOMC"
        if ! start-stop-daemon --start --oknodo --quiet --make-pidfile --background \
            --pidfile "$PIDFILE" --chuid "$DAEMONUSER" \
            --exec /bin/bash -- -c "$DAEMON $DAEMON_ARGS >> $LOGFILE 2>&1"
        then
            log_end_msg 1
        else
            log_end_msg 0
        fi
        ;;
    stop)
        log_daemon_msg "Stopping $DESC" "$NAME"
        if start-stop-daemon --stop --retry 30 --oknodo --quiet \
            --pidfile "$PIDFILE" --exec "$DAEMON"
        then
            kill `ps ax | grep "local/bin/plackup" | grep -v "grep" | awk '{print $1}'`
            rm -f "$PIDFILE"
            log_end_msg 0
        else
            log_end_msg 1
        fi
        ;;
    restart)
        $0 stop
        $0 start
        ;;
    status)
        status_of_proc -p "$PIDFILE" "$DAEMON" "$NAME" && exit 0 || exit $?
        ;;
    *)
        echo "Usage: $SCRIPTNAME {start|stop|restart|status}" >&2
        exit 1
        ;;
esac

exit 0
