# containers/cron-refresh-stats/crontab
# m h dom mon dow user    command
46  3 *   *   sun root    /opt/mediacloud/bin/cron_refresh_stats.pl 1> /dev/stdout 2> /dev/stderr
