# containers/cron-refresh-stats/crontab
# m h dom mon dow user    command
46  3 *   *   sun root    /cron_refresh_stats.pl 1> /proc/1/fd/1 2> /proc/1/fd/2
