# containers/cron-rescraping-changes/crontab
# m h dom mon dow user    command
0   3 *   *   *   root    /opt/mediacloud/bin/cron_rescraping_changes.py 1> /dev/stdout 2> /dev/stderr
