# supervisor/supervisord.conf
;
; Media Cloud configuration file for Supervisor (http://supervisord.org/)
;
; To add your own services to Supervisor, create
; "supervisord.user.service_name.conf" and add your additional configuration
; properties there.
;

[inet_http_server]
port = 127.0.0.1:4398
username = supervisord
password = qHujfp7n4J

[supervisord]
logfile = %(here)s/../data/supervisor_logs/supervisord.log  ; main log file (default $CWD/supervisord.log)
pidfile = %(here)s/../data/supervisor_logs/supervisord.pid     ; supervisord pidfile (default supervisord.pid)
# DO NOT set "childlogdir" here because it is set by the supervisord.sh script
;childlogdir = %(here)s/../data/supervisor_logs/            ; ('AUTO' child log dir, default $TEMP)
directory = %(here)s/../                                    ; chroot to Media Cloud's root directory
nocleanup = true                                            ; don't remove process logs after stopping the process

; The below section must remain in the config file for RPC (supervisorctl /
; web interface) to work, additional interfaces may be added by defining them
; in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix://%(here)s/supervisord.sock                ; use a unix:// URL for a unix socket

[group:mc]
programs = crawler,process_story_subsets,update_processed_stories,extract_and_vector,author_extractor,update_processed_stories
;priority=999

[program:crawler]
command = bash ./script/run_with_carton.sh ./script/mediawords_crawl.pl 10
directory = %(here)s/../
;stopwaitsecs = 20
autorestart = false
autostart = true

[program:process_story_subsets]
command = bash ./script/run_with_carton.sh ./script/mediawords_process_story_subsets.pl
directory = %(here)s/../
autorestart = false
autostart = true

[program:update_processed_stories]
command = bash ./script/run_with_carton.sh ./script/mediawords_update_processed_stories.pl -d
directory = %(here)s/../
autorestart = false
autostart = true

[program:extract_and_vector]
command = bash ./script/run_with_carton.sh ./script/mediawords_extract_and_vector.pl 10
directory = %(here)s/../
;stopwaitsecs = 20
autorestart = false
autostart = true

[program:author_extractor]
command = bash ./script/run_with_carton.sh ./script/mediawords_extract_author.pl 1
directory = %(here)s/../
;stopwaitsecs = 20
autorestart = false
autostart = false

[program:update_aggregate_vectors]
command = bash ./script/run_with_carton.sh ./script/mediawords_update_aggregate_vectors.pl -d
directory = %(here)s/../
stopwaitsecs = 20
autorestart = true
autostart = false

; Gearman (p5-Gearman-JobScheduler) functions
[group:gearman]
programs = add_default_feeds,search_stories
;priority=999

[program:add_default_feeds]
command = bash ./script/run_with_carton.sh local/bin/gjs_worker.pl lib/MediaWords/GearmanFunction/AddDefaultFeeds.pm
directory = %(here)s/../
autorestart = true
autostart = true
numprocs = 2
process_name = %(program_name)s_%(process_num)02d

[program:search_stories]
command = bash ./script/run_with_carton.sh local/bin/gjs_worker.pl lib/MediaWords/GearmanFunction/SearchStories.pm
directory = %(here)s/../
autorestart = true
autostart = true
numprocs = 2
process_name = %(program_name)s_%(process_num)02d

[program:mine_controversy]
command = bash ./script/run_with_carton.sh local/bin/gjs_worker.pl lib/MediaWords/GearmanFunction/CM/MineControversy.pm
directory = %(here)s/../
autorestart = true
autostart = true
numprocs = 1
process_name = %(program_name)s_%(process_num)02d

[program:dump_controversy]
command = bash ./script/run_with_carton.sh local/bin/gjs_worker.pl lib/MediaWords/GearmanFunction/CM/DumpControversy.pm
directory = %(here)s/../
autorestart = true
autostart = true
numprocs = 1
process_name = %(program_name)s_%(process_num)02d

; Gearman (p5-Gearman-JobScheduler) functions
[group:solr_related]
programs = solr
;priority=999

[program:solr]
command = bash ./script/solr/solr_start.sh
directory = %(here)s/../
stopwaitsecs = 20
autorestart = true
autostart = false

[group:thrift]
programs=solr_facets_wrapper_server

[program:solr_facets_wrapper_server]
command = python  python_scripts/solr_facets_wrapper_server.py
directory = %(here)s/../
autorestart = false
autostart = false


; Include custom user's configuration
[include]
files = supervisord.user.*.conf
