# apps/cron-refresh-stats/crontab
# Don't send no email
MAILTO=""

# m h dom mon dow user          command
46  3 *   *   sun mediacloud    /opt/mediacloud/bin/cron_refresh_stats.pl 1> /dev/stdout 2> /dev/stderr
