# containers/cron-rescrape-due-media/crontab
# m h dom mon dow user    command
0   3 *   *   *   root    /opt/mediacloud/bin/cron_rescrape_due_media.pl 1> /dev/stdout 2> /dev/stderr
