# apps/cron-print-long-running-job-states/crontab
# m h  dom mon dow user         command
6   19 *   *   *   mediacloud   /opt/mediacloud/bin/cron_print_long_running_job_states.py 1> /dev/stdout 2> /dev/stderr
