# apps/cron-set-media-subject-country/crontab
# m h dom mon dow user    command
16  2 *   *   *   root    /opt/mediacloud/bin/cron_set_media_subject_country.pl 1> /dev/stdout 2> /dev/stderr
